# chef/cookbooks/baseconfig/files/default/rc.local
#!/bin/sh
sleep 10
/usr/local/bin/uwsgi --ini /home/vagrant/project/uwsgi.ini --daemonize /var/log/blogproject.log
